# scripts/TaskScheduler_Scheduler.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler definition for "MemeDNA Pipeline Scheduler".

  This hosts the long-running APScheduler loop (memedna.pipeline.scheduler),
  which ticks every PIPELINE_INTERVAL_MINUTES and runs the full pipeline.

  Import:
    schtasks /Create /TN "MemeDNA Pipeline Scheduler" /XML "scripts\TaskScheduler_Scheduler.xml" /F
  Start now (without waiting for logon):
    schtasks /Run /TN "MemeDNA Pipeline Scheduler"
  Stop:
    schtasks /End /TN "MemeDNA Pipeline Scheduler"

  Edit the `UserId` and paths below before importing on a new machine.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Keep MemeDNA's APScheduler alive — runs the ingest→cluster→analyse pipeline every 5 minutes.</Description>
    <Author>MemeDNA</Author>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
    <BootTrigger>
      <Enabled>true</Enabled>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>CHANGE_ME_USERNAME</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>999</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "C:\Users\user\Desktop\MemeDNA\scripts\run_scheduler.ps1"</Arguments>
      <WorkingDirectory>C:\Users\user\Desktop\MemeDNA</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
